# apps/gateway_srv/sup_gateway.conf
[program:gateway]
directory=/rel/apps/gateway
command=/rel/apps/gateway/gateway_srv
redirect_stderr=true
stdout_logfile=/rel/apps/gateway/log/access.log
