# apps/identity/sup_identity.conf
[program:identity_prod]
directory=/rel/apps/identity_prod
command=/rel/apps/identity_prod/identity_srv
redirect_stderr=true
stdout_logfile=/rel/apps/identity_prod/log/access.log