# __deprecated/gateway_http/sup_gateway.conf
[program:gateway]
directory=/mega/apps/gateway_http
command=/mega/apps/gateway_http/gateway_http
redirect_stderr=true
stdout_logfile=/mega/apps/gateway_http/log/access.log
